# rockchip/common/overlays/rootfs/usb-gadget/lib/udev/rules.d/61-usbdevice.rules
SUBSYSTEM=="android_usb", ACTION=="change", RUN+="/usr/bin/usbdevice update"
SUBSYSTEM=="udc", RUN+="/usr/bin/usbdevice update"
